# iiif/supervisord.conf
[supervisord]
nodaemon=true

[program:iiifserver_license]
command=/usr/share/cgi-bin/iiifserver.fcgi --activate
environment = VERBOSITY="0"
user=user
autostart=true
autorestart=false
priority=1
startsecs=0

[fcgi-program:iiifserver]
command=/usr/share/cgi-bin/iiifserver.fcgi
socket=tcp://localhost:9000
process_name=%(program_name)s_%(process_num)02d
user=user
numprocs=8
autostart=true
autorestart=true
priority=999
stopsignal=TERM
environment = LOGFILE="/tmp/iiifserver.log",
        VERBOSITY="2",
        MEMCACHED_SERVERS="localhost",
        MEMCACHED_TIMEOUT="0",
        MAX_CVT="1000",
        FILESYSTEM_PREFIX="/data/",
        FILENAME_EXTENSION=".jp2",
        VIEWER="/viewer.html"

[program:memcached]
command=/usr/bin/memcached -p 11211 -u memcached -m 1024 -c 1024 -t 4
autostart=true
autorestart=true
user=user
priority=100
redirect_stderr=true


[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true

[program:dropzone]
command=/usr/bin/dropzone.py
autostart=true
autorestart=true
logfile=/tmp/drop_supervisord.log
user=user
#[program:start_copy]
#command=/bin/cp /demo.jp2 /data/
#autostart=true
#startsecs=0

