# iiif/iiifserver/supervisord.conf
[supervisord]
nodaemon=true

[program:iiifserver_license]
command=/usr/share/cgi-bin/iiifserver.fcgi --activate
environment = VERBOSITY="0"
autostart=true
autorestart=false
priority=1
startsecs=0
stderr_logfile=/var/log/nginx/license_stderr.log

[program:memcached]
command=/usr/bin/memcached -p 11211 -u memcached -m 1024 -c 1024 -t 4
autostart=true
autorestart=true
priority=100
redirect_stderr=true
user=user
stderr_logfile=/var/log/nginx/memcached_err.log

[fcgi-program:iiifserver]
command=/usr/share/cgi-bin/iiifserver.fcgi
socket=tcp://localhost:9000
process_name=%(program_name)s_%(process_num)02d
user=user
numprocs=8
autostart=true
autorestart=true
priority=999
stopsignal=TERM
environment = LOGFILE="/var/log/nginx/iiifserver.log",
        VERBOSITY="10",
        MEMCACHED_SERVERS="localhost",
        MEMCACHED_TIMEOUT="0",
        MAX_IMAGE_CACHE_SIZE=%(ENV_MAX_IMAGE_CACHE_SIZE)s,
        JPEG_QUALITY=%(ENV_JPEG_QUALITY)s,
        MAX_CVT=%(ENV_MAX_CVT)s,
        CORS="*",
        FILESYSTEM_PREFIX="/data/",
        FILENAME_EXTENSION=".jp2",
        VIEWER="/viewer.html",


[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx/nginx_stderr.log

[program:dropzone]
command=/usr/bin/dropzone.py
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx/dropzone_err.log
user=user


